# hydroweb/hydroweb.socket
[Unit]
Description=Hydroweb Gunicorn socket

[Socket]
ListenStream=/run/hydroweb.sock

[Install]
WantedBy=sockets.target
